# etc/network/if-up.d/powerdns
#!/bin/bash

echo "nameserver 127.0.0.1" | cat - /etc/resolv.conf > /tmp/.inkspot-resolv && mv /tmp/.inkspot-resolv /etc/resolv.conf
